# Sources/iAura/Resources/com.user.iaura.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- LaunchAgent template for iAura daemon. -->
<dict>
    <key>Label</key>
    <string>com.user.iaura</string>
    <key>ProgramArguments</key>
    <array>
        <string>BINARY_PATH_PLACEHOLDER</string>
        <string>serve</string>
    </array>
    <key>WorkingDirectory</key>
    <string>WORKING_DIR_PLACEHOLDER</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
